# deploy/systemd/biogestor-backend.service
# Servicio systemd para backend ASGI de Biogestor
# Ejecuta Daphne (ASGI) con el entorno virtual del proyecto

[Unit]
Description=Biogestor Backend (Daphne ASGI)
After=network.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/biogestor/backend
# Carga variables de entorno del proyecto
EnvironmentFile=/srv/biogestor/backend/BGProject/.env
# Activa el venv y ejecuta Daphne
ExecStart=/srv/biogestor/.venv/bin/daphne -b 127.0.0.1 -p 8000 BGProject.asgi:application
Restart=always
RestartSec=3
# Limitar recursos (opcional)
# MemoryMax=512M

# Asegura acceso a staticfiles/media
ReadWritePaths=/srv/biogestor/backend/media

[Install]
WantedBy=multi-user.target
